# ansible/jetty9-install/roles/oskari-server/templates/error-notifier.service.j2
[Unit]
Description=OnFailure for %i

[Service]
Type=oneshot
ExecStart={{ oskari_dir }}/error_notify.sh %i
